# apps/ws/supervisor-ws.conf
[program:guard51-ws]
command=node /www/wwwroot/guard51/apps/ws/server.js
directory=/www/wwwroot/guard51/apps/ws
environment=WS_PORT=8089
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/var/log/guard51-ws.log
stderr_logfile=/var/log/guard51-ws-error.log
