# apps/api/config/supervisor-worker.conf
[program:guard51-worker]
command=php /www/wwwroot/guard51/apps/api/bin/worker.php
directory=/www/wwwroot/guard51/apps/api
autostart=true
autorestart=true
startsecs=3
startretries=10
user=www-data
stdout_logfile=/var/log/guard51-worker.log
stderr_logfile=/var/log/guard51-worker-error.log
